# deployments/supervisord.d/fivc-club-web.ini
[program:fivc-club-web]
user=root
directory=/opt/app/fivc-club
command=/opt/app/fivc-club/venv/bin/gunicorn -w 2 -t 180 -k gevent --pythonpath=. -b 0.0.0.0:8000 launchers.wsgi
environment=PATH="/opt/app/fivc-club/venv/bin",PYTHONPATH="/opt/app/fivc-club"
stdout_logfile=/var/log/fivc-club/fivc-club-web.log
stderr_logfile=/var/log/fivc-club/fivc-club-web-error.log
autostart=true
autorestart=true
startsecs=10
